# tester.sv.conf
[program:tester]
command=java -agentlib:TakipiAgent -jar /scala-boom.jar
